# docker/supervisor.d/sqlitebackup.ini
[program:sqlbackup]
command=/usr/local/bin/sqlitebackup.sh -r
autostart=true
autorestart=true
stderr_logfile=/var/log/sqlbackup.err.log
stdout_logfile=/var/log/sqlbkacup.out.log
