# docker/supervisor.d/sqlitebackup.ini
[program:sqlbackup]
command=/usr/local/bin/sqlitebackup.sh -r
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
